# ansible/roles/gitea/templates/gitea.service.j2
[Unit]
Description=Gitea (Git with a cup of tea)
After=syslog.target
After=network.target
After=mysql.service

[Service]
RestartSec=2s
Type=simple
User={{ gitea_user }}
Group={{ gitea_user }}
WorkingDirectory={{ gitea_work_dir }}
ExecStart=/usr/local/bin/gitea web -c {{ gitea_config_dir }}/app.ini
Restart=always
Environment=USER={{ gitea_user }} HOME={{ gitea_home }} GITEA_WORK_DIR={{ gitea_work_dir }}

[Install]
WantedBy=multi-user.target
